# etc/profile.d/environment.sh
export DOMAIN=
export PUID=
export PGID=
export ACME_EMAIL=
export CF_EMAIL=
export CF_API_KEY=
export NAMESERVER_1=
export NAMESERVER_2=
export OAUTH_CLIENT_ID=
export OAUTH_CLIENT_SECRET=
export OAUTH_SECRET=
export OAUTH_WHITELIST=
export MYSQL_PASSWORD=
export RCLONE_CRON=
export RCLONE_CRON_ABORT=
export STEAM_USERNAME=
export STEAM_PASSWORD=
export MINECRAFT_VERSION=
export MINECRAFT_MEMORY=
export APPDATA=
export APPDATA_WEB=
export APPDATA_GAMES=
